# s6/services.d/code-server/finish
#!/usr/bin/execlineb -S0
# =============================================================================
# s6-overlay finish script for code-server
# Handles service termination
# =============================================================================

# Log the exit code
foreground {
    if { s6-echo "[services.d] code-server exited with code ${1}" }
}

# Exit with the same code
exit ${1}
