# s6/services.d/code-server/run
#!/usr/bin/with-contenv bash
# =============================================================================
# s6-overlay service script for code-server
# Starts VS Code Server on port 8888 with Kubeflow compatibility
# =============================================================================

set -e

echo "[services.d] Starting code-server service..."

# Change to project directory
cd /home/${NB_USER}/project

# Start code-server with proper configuration
# - Bind to 0.0.0.0:8888 (Kubeflow requirement)
# - No authentication (Kubeflow handles auth)
# - Disable telemetry for privacy
# - Use configuration file for additional settings
exec s6-setuidgid ${NB_USER} \
    code-server \
    --bind-addr 0.0.0.0:8888 \
    --auth none \
    --disable-telemetry \
    --disable-update-check \
    --config /etc/code-server/config.yaml \
    --extensions-dir /home/${NB_USER}/.local/share/code-server/extensions \
    --user-data-dir /home/${NB_USER}/.local/share/code-server \
    /home/${NB_USER}/project
